# failgather.cron
0 * * * * pi cd ~/failgather ; ./gather.sh
